# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:shrpx]
command=/usr/local/bin/shrpx

[program:squid3]
command=/usr/sbin/squid3 -N

[program:cronjob]
command=/usr/sbin/cron -f

